# clipmenu.timer
[Unit]
Description=Run clipmenud every tenth of a second.

[Timer]
OnUnitInactiveSec=500us
AccuracySec=300us

[Install]
WantedBy=timers.target
